# supervisor.ini
[inet_http_server]
port = :9001

[program:Correspond]
directory = ./
command = ./Correspond.exe

[program:LogonServer]
directory = ./
command = ./LogonServer.exe

[program:ControlServer]
directory = ./
command = ./ControlServer.exe /auto:1

[program:GameServer]
directory = ./
command = ./GameServer.exe /ServerID:28
